# scripts/com.skillnav.failover-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.skillnav.failover-check</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/node</string>
		<string>/Users/apple/WeChatProjects/skillnav/scripts/failover-check.mjs</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/apple/WeChatProjects/skillnav</string>
	<key>StartInterval</key>
	<integer>3600</integer>
	<key>StandardOutPath</key>
	<string>/tmp/skillnav-failover.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/skillnav-failover.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin</string>
	</dict>
</dict>
</plist>
